# roles/cbm-machine/files/set-hostname.service
[Unit]
[Service]
Type=oneshot
ExecStart=/usr/bin/hostnamectl set-hostname $HOSTNAME
[Install]
WantedBy=multi-user.target